# systemd/meshcore-monitor.service
[Unit]
Description=MeshCore Repeater Pi Monitor
After=network.target

[Service]
Type=simple
User=meshcoremon
WorkingDirectory=/opt/meshcore-monitor
ExecStart=/opt/meshcore-monitor/venv/bin/python app.py
Restart=always
RestartSec=10
# Configuration is loaded from .env file in WorkingDirectory (see .env.example)

[Install]
WantedBy=multi-user.target
